# roles/tg/files/pa.service
[Unit]
Description=PA service

[Service]
WorkingDirectory=/var/pa
ExecStart=/usr/bin/python3 /var/pa/app.py

[Install]
WantedBy=default.target
